# linux_service/coderco-app.service
[Unit]
Description=CoderCo App Service
After=network.target

[Service]
User=appuser
WorkingDirectory=/opt/coderco-app
ExecStart=/usr/bin/python3 /opt/coderco-app/server.py
EnvironmentFile=/opt/coderco-app/.env
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target